# etc/crontab
0 1 * * * /opt/gettab-server/scripts/get_wp
*/10 * * * * /opt/wallpapers/update.sh
